# internal/policies/mount/testdata/TestCreateUnits/error_when_writing_multiple_units_with_badly_formatted_among_them/golden/adsys-correct:--format.com-value.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for correct://format.com/value
After=network-online.target

[Mount]
What=correct://format.com/value
Where=/mnt/adsys/format.com/value
Options=nofail,

[Install]
WantedBy=default.target
