# internal/policies/mount/testdata/golden/TestApplyPolicy/system,_mount_units_are_updated_on_refreshing_policy_with_an_entry_with_multiple_values/etc/systemd/system/adsys-protocol-domain.com-mountpath2.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for protocol://domain.com/mountpath2
After=network-online.target

[Mount]
What=/domain.com/mountpath2
Where=/adsys/protocol/domain.com/mountpath2
Type=protocol
Options=default

[Install]
WantedBy=default.target
