# internal/policies/mount/testdata/golden/TestApplyPolicy/error_when_trying_to_update_policy_with_badly_formatted_entry/etc/systemd/system/adsys-protocol-domain.com-mountpath.mount
# This template defines the basic structure of a mount unit generated by ADSys for system mounts.
[Unit]
Description=ADSys mount for protocol://domain.com/mountpath
After=network-online.target
Requires=network-online.target

[Mount]
What=/domain.com/mountpath
Where=/adsys/protocol/domain.com/mountpath
Type=protocol
Options=defaults
# This option prevents hangs on shutdown due to an unreachable network share.
LazyUnmount=true
TimeoutSec=30

[Install]
WantedBy=default.target
