# internal/policies/mount/testdata/TestCreateUnits/write_anonymous_unit_with_correct_options/golden/adsys-protocol:--anon.com-anon-mount.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for [anonymous]protocol://anon.com/anon/mount
After=network-online.target

[Mount]
What=protocol://anon.com/anon/mount
Where=/mnt/adsys/anon.com/anon/mount
Options=nofail,users

[Install]
WantedBy=default.target
