# internal/policies/mount/testdata/golden/TestCreateUnits/write_multiple_units/adsys-nfs-yetanotherdomain.com-mount_path-mount-path.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for nfs://yetanotherdomain.com/mount_path/mount/path
After=network-online.target

[Mount]
What=yetanotherdomain.com:/mount_path/mount/path
Where=/adsys/nfs/yetanotherdomain.com/mount_path/mount/path
Type=nfs
Options=default

[Install]
WantedBy=default.target
