# internal/policies/mount/testdata/TestApplyPolicy/system,_successfully_apply_policy_with_anonymous_values/etc-systemd-system/adsys-smb:--shady.com-shady-mount.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for [anonymous]smb://shady.com/shady/mount
After=network-online.target

[Mount]
What=smb://shady.com/shady/mount
Where=/mnt/adsys/shady.com/shady/mount
Options=nofail,users

[Install]
WantedBy=default.target
