# internal/policies/mount/testdata/TestApplyPolicy/system,_successfully_apply_policy_with_anonymous_values/etc-systemd-system/adsys-nfs:--stealthydomain.com-mount-stealthy.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for [anonymous]nfs://stealthydomain.com/mount/stealthy
After=network-online.target

[Mount]
What=nfs://stealthydomain.com/mount/stealthy
Where=/mnt/adsys/stealthydomain.com/mount/stealthy
Options=nofail,users

[Install]
WantedBy=default.target
