# internal/policies/mount/testdata/TestApplyPolicy/system,_successfully_apply_policy_for_entry_with_repeatead_values/etc-systemd-system/adsys-rpt:--repeated.com-repeatedmount.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for rpt://repeated.com/repeatedmount
After=network-online.target

[Mount]
What=rpt://repeated.com/repeatedmount
Where=/mnt/adsys/repeated.com/repeatedmount
Options=nofail,

[Install]
WantedBy=default.target
