# internal/policies/mount/testdata/TestApplyPolicy/system,_successfully_apply_policy_for_entry_with_repeatead_values/etc-systemd-system/adsys-nfs:--anotherone.com-mnt.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for nfs://anotherone.com/mnt
After=network-online.target

[Mount]
What=nfs://anotherone.com/mnt
Where=/mnt/adsys/anotherone.com/mnt
Options=nofail,

[Install]
WantedBy=default.target
