# internal/policies/mount/testdata/TestApplyPolicy/system,_successfully_apply_policy_for_entry_with_repeated_values/etc/systemd/system/adsys-smb:--single.com-mnt.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for smb://single.com/mnt
After=network-online.target

[Mount]
What=smb://single.com/mnt
Where=/mnt/adsys/single.com/mnt

[Install]
WantedBy=default.target
