# internal/policies/mount/testdata/TestApplyPolicy/system,_successfully_apply_policy_with_kerberos_tagged_values/etc/systemd/system/adsys-smb:--authenticated.com-authenticated-mount.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for [krb5]smb://authenticated.com/authenticated/mount
After=network-online.target

[Mount]
What=smb://authenticated.com/authenticated/mount
Where=/mnt/adsys/authenticated.com/authenticated/mount

[Install]
WantedBy=default.target
