# internal/policies/mount/testdata/TestApplyPolicy/system,_mount_units_are_added_on_refreshing_policy_with_slightly_different_values/etc/systemd/system/adsys-nfs:--yetanotherdomain.com-mount\x2dpath-mount-path.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for nfs://yetanotherdomain.com/mount-path/mount/path
After=network-online.target

[Mount]
What=nfs://yetanotherdomain.com/mount-path/mount/path
Where=/mnt/adsys/yetanotherdomain.com/mount-path/mount/path
Options=nofail,

[Install]
WantedBy=default.target
