# internal/policies/mount/testdata/golden/TestApplyPolicy/system,_successfully_apply_policy_for_entry_with_repeated_values/etc/systemd/system/adsys-rpt-repeated.com-repeatedmount.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for rpt://repeated.com/repeatedmount
After=network-online.target

[Mount]
What=/repeated.com/repeatedmount
Where=/adsys/rpt/repeated.com/repeatedmount
Type=rpt
Options=default

[Install]
WantedBy=default.target
