# internal/policies/mount/testdata/golden/TestApplyPolicy/system,_successfully_apply_policy_for_entry_with_repeated_values/etc/systemd/system/adsys-nfs-anotherone.com-mnt.mount
# This template defines the basic structure of a mount unit generated by ADSys for system mounts.
[Unit]
Description=ADSys mount for nfs://anotherone.com/mnt
After=network-online.target
Requires=network-online.target

[Mount]
What=anotherone.com:/mnt
Where=/adsys/nfs/anotherone.com/mnt
Type=nfs
Options=defaults
# This option prevents hangs on shutdown due to an unreachable network share.
LazyUnmount=true
TimeoutSec=30

[Install]
WantedBy=default.target
