# internal/policies/mount/testdata/golden/TestApplyPolicy/system,_successfully_apply_policy_for_entry_with_repeated_values/etc/systemd/system/adsys-nfs-anotherone.com-mnt.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for nfs://anotherone.com/mnt
After=network-online.target

[Mount]
What=anotherone.com:/mnt
Where=/adsys/nfs/anotherone.com/mnt
Type=nfs
Options=default

[Install]
WantedBy=default.target
