# internal/policies/mount/testdata/golden/TestApplyPolicy/system,_successfully_apply_policy_for_entry_with_repeated_values/etc/systemd/system/adsys-cifs-single.com-mnt.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for smb://single.com/mnt
After=network-online.target

[Mount]
What=//single.com/mnt
Where=/adsys/cifs/single.com/mnt
Type=cifs
Options=default

[Install]
WantedBy=default.target
