# systemd/adsys-gpo-refresh.service
[Unit]
Description=Refresh ADSys GPO for machine and users

[Service]
Type=oneshot
ExecStart=/sbin/adsysctl update --all
