# systemd/adsys-gpo-refresh.timer
[Unit]
Description=Refresh ADSys GPO for machine and users

[Timer]
OnBootSec=30min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target