# systemd/adsysd.socket
[Unit]
Description=Socket activation for ADSys daemon

[Socket]
ListenStream=/run/adsysd.sock
FlushPending=yes

[Install]
WantedBy=sockets.target
