# internal/policies/mount/testdata/golden/TestCreateUnits/write_krb5_tagged_unit/adsys-protocol-kerberos.com-auth_mount.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for [krb5]protocol://kerberos.com/auth_mount
After=network-online.target

[Mount]
What=/kerberos.com/auth_mount
Where=/adsys/protocol/kerberos.com/auth_mount
Type=protocol
Options=sec=krb5i

[Install]
WantedBy=default.target
