# internal/policies/mount/testdata/TestCreateUnits/write_single_unit/golden/adsys-protocol:--domain.com-mountpath.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for protocol://domain.com/mountpath
After=network-online.target

[Mount]
What=protocol://domain.com/mountpath
Where=/mnt/adsys/domain.com/mountpath
Options=nofail,

[Install]
WantedBy=default.target
