# internal/policies/mount/testdata/golden/TestApplyPolicy/system,_successfully_apply_policy_for_entry_with_multiple_values/etc/systemd/system/adsys-protocol:--domain.com-mountpath2.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for protocol://domain.com/mountpath2
After=network-online.target

[Mount]
What=protocol://domain.com/mountpath2
Where=/mnt/adsys/domain.com/mountpath2

[Install]
WantedBy=default.target
