# internal/policies/mount/testdata/golden/TestApplyPolicy/system,_successfully_apply_policy_for_entry_with_multiple_values/etc/systemd/system/adsys-cifs-otherdomain.com-mount-path.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for smb://otherdomain.com/mount/path
After=network-online.target

[Mount]
What=//otherdomain.com/mount/path
Where=/adsys/cifs/otherdomain.com/mount/path
Type=cifs
Options=default

[Install]
WantedBy=default.target
