# internal/policies/testdata/golden/applypolicy/succeed/etc/systemd/system/adsys-fuse-example.com-ftp_share.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for ftp://example.com/ftp_share
After=network-online.target

[Mount]
What=example.com
Where=/adsys/fuse/example.com/ftp_share
Type=fuse
Options=default

[Install]
WantedBy=default.target
