# internal/policies/testdata/golden/applypolicy/succeed/etc/systemd/system/adsys-smb:--example.com-smb_share.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for smb://example.com/smb_share
After=network-online.target

[Mount]
What=smb://example.com/smb_share
Where=/mnt/adsys/example.com/smb_share

[Install]
WantedBy=default.target
