# internal/policies/testdata/golden/applypolicy/succeed/etc/systemd/system/adsys-nfs-example.com-nfs_share.mount
# This template defines the basic structure of a mount unit generated by ADSys for system mounts.
[Unit]
Description= ADSys mount for nfs://example.com/nfs_share
After=network-online.target

[Mount]
What=example.com:/nfs_share
Where=/adsys/nfs/example.com/nfs_share
Type=nfs
Options=defaults
# This option prevents hangs on shutdown due to an unreachable network share.
LazyUnmount=true

[Install]
WantedBy=default.target
