# internal/policies/testdata/golden/applypolicy/succeed/etc/systemd/system/adsys-ftp:--example.com-ftp_share.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for ftp://example.com/ftp_share
After=network-online.target

[Mount]
What=ftp://example.com/ftp_share
Where=/mnt/adsys/example.com/ftp_share

[Install]
WantedBy=default.target
