# internal/policies/testdata/golden/applypolicy/succeed/etc/systemd/system/adsys-cifs-example.com-smb_share.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for smb://example.com/smb_share
After=network-online.target

[Mount]
What=//example.com/smb_share
Where=/adsys/cifs/example.com/smb_share
Type=cifs
Options=default

[Install]
WantedBy=default.target
