# internal/policies/testdata/golden/applypolicy/succeed/etc/systemd/system/adsys-cifs-example.com-smb_share.mount
# This template defines the basic structure of a mount unit generated by ADSys for system mounts.
[Unit]
Description= ADSys mount for smb://example.com/smb_share
After=network-online.target
Requires=network-online.target

[Mount]
What=//example.com/smb_share
Where=/adsys/cifs/example.com/smb_share
Type=cifs
Options=defaults
# This option prevents hangs on shutdown due to an unreachable network share.
LazyUnmount=true
TimeoutSec=30

[Install]
WantedBy=default.target
