# internal/policies/testdata/golden/applypolicy/succeed/etc/systemd/system/adsys-nfs:--example.com-nfs_share.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for nfs://example.com/nfs_share
After=network-online.target

[Mount]
What=nfs://example.com/nfs_share
Where=/mnt/adsys/example.com/nfs_share

[Install]
WantedBy=default.target
