# internal/policies/mount/testdata/TestApplyPolicy/error_when_stopping_units_for_clean_up_and_systemctl_fails/etc-systemd-system/adsys-protocol:--otherdomain.com-mountpath2.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for protocol://otherdomain.com/mountpath2
After=network-online.target

[Mount]
What=protocol://otherdomain.com/mountpath2
Where=/mnt/adsys/otherdomain.com/mountpath2
Options=nofail,

[Install]
WantedBy=default.target
