# internal/policies/mount/adsys-mount-template.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for %s
After=network-online.target

[Mount]
What=%s
Where=/mnt/adsys/%s

[Install]
WantedBy=default.target
