# internal/policies/mount/testdata/golden/TestApplyPolicy/system,_successfully_apply_policy_trimming_sequential_linebreaks/etc/systemd/system/adsys-protocol-domain.com-mounpath.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for protocol://domain.com/mounpath
After=network-online.target

[Mount]
What=/domain.com/mounpath
Where=/adsys/protocol/domain.com/mounpath
Type=protocol
Options=default

[Install]
WantedBy=default.target
