# internal/policies/mount/testdata/golden/TestApplyPolicy/system,_successfully_apply_policy_with_kerberos_tagged_values/etc/systemd/system/adsys-cifs-authenticated.com-authenticated-mount.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for [krb5]smb://authenticated.com/authenticated/mount
After=network-online.target

[Mount]
What=//authenticated.com/authenticated/mount
Where=/adsys/cifs/authenticated.com/authenticated/mount
Type=cifs
Options=sec=krb5i

[Install]
WantedBy=default.target
