# internal/policies/mount/testdata/golden/TestApplyPolicy/system,_successfully_apply_policy_with_kerberos_tagged_values/etc/systemd/system/adsys-nfs-krb_domain.com-mount-krb_path.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for [krb5]nfs://krb_domain.com/mount/krb_path
After=network-online.target

[Mount]
What=krb_domain.com:/mount/krb_path
Where=/adsys/nfs/krb_domain.com/mount/krb_path
Type=nfs
Options=sec=krb5i

[Install]
WantedBy=default.target
