# internal/policies/mount/testdata/TestCreateUnits/write_krb5_tagged_unit/golden/adsys-protocol:--kerberos.com-auth_mount.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for [krb5]protocol://kerberos.com/auth_mount
After=network-online.target

[Mount]
What=protocol://kerberos.com/auth_mount
Where=/mnt/adsys/kerberos.com/auth_mount

[Install]
WantedBy=default.target
