# internal/policies/mount/testdata/golden/TestApplyPolicy/system,_successfully_apply_policy_trimming_whitespaces/etc/systemd/system/adsys-nfs:--yetanotherdomain.com-path-mount.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for nfs://yetanotherdomain.com/path/mount
After=network-online.target

[Mount]
What=nfs://yetanotherdomain.com/path/mount
Where=/mnt/adsys/yetanotherdomain.com/path/mount

[Install]
WantedBy=default.target
