# internal/policies/mount/testdata/golden/TestApplyPolicy/system,_successfully_apply_policy_trimming_whitespaces/etc/systemd/system/adsys-protocol-domain.com-mountpath.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for protocol://domain.com/mountpath
After=network-online.target

[Mount]
What=/domain.com/mountpath
Where=/adsys/protocol/domain.com/mountpath
Type=protocol
Options=default

[Install]
WantedBy=default.target
