# internal/policies/mount/testdata/TestApplyPolicy/system,_mount_units_are_updated_on_refreshing_policy_with_an_entry_with_multiple_values/etc/systemd/system/adsys-smb:--otherdomain.com-mount-path.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for smb://otherdomain.com/mount/path
After=network-online.target

[Mount]
What=smb://otherdomain.com/mount/path
Where=/mnt/adsys/otherdomain.com/mount/path
Options=nofail,

[Install]
WantedBy=default.target
