# internal/policies/mount/testdata/TestApplyPolicy/system,_mount_units_are_added_on_refreshing_policy_with_some_matching_values/etc/systemd/system/adsys-ftp:--completelydifferent.com-different-path.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for ftp://completelydifferent.com/different/path
After=network-online.target

[Mount]
What=ftp://completelydifferent.com/different/path
Where=/mnt/adsys/completelydifferent.com/different/path

[Install]
WantedBy=default.target
