# internal/policies/mount/testdata/TestApplyPolicy/system,_mount_units_are_removed_on_refreshing_policy_with_an_entry_with_different_values/etc/systemd/system/adsys-nfs:--domain.com-mount\x2dpath.mount
# This template defines the basic structure of a mount unit generated by adsys for system mounts.
[Unit]
Description=Mount for nfs://domain.com/mount-path
After=network-online.target

[Mount]
What=nfs://domain.com/mount-path
Where=/mnt/adsys/domain.com/mount-path
Options=nofail,

[Install]
WantedBy=default.target
